# crontab
0 3 * * * /app/mongodb-backup.sh >> /proc/1/fd/1 2>> /proc/1/fd/2